# ops/systemd/cert-expiry-alert.service
[Unit]
Description=Studytracker TLS certificate expiry alert
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/apps/studytracker
ExecStart=/usr/bin/python3 /root/apps/studytracker/scripts/cert_expiry_alert.py
